# supervisor-app.conf
[program:app-uwsgi]
command = /usr/local/bin/uwsgi --ini /home/docker/code/jpsp/uwsgi.ini
stopsignal=QUIT
autostart=true
autorestart=true